# kgitlab.service.in
[Unit]
Description=GitLab System Hook Listener to Enable Kerberos
After=network.target

[Service]
ExecStart=@bindir@/kgitlab server --config @sysconfdir@/kgitlab/config.yaml

[Install]
WantedBy=multi-user.target
